# com.user.photo-memories.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.photo-memories</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/wuchanghong/projects/ai-photo-selector/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wuchanghong/projects/ai-photo-selector</string>

    <key>StandardOutPath</key>
    <string>/Users/wuchanghong/projects/ai-photo-selector/logs/photo-memories.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wuchanghong/projects/ai-photo-selector/logs/photo-memories.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
